# data/threathunts/hunts.json
{
  "engine": "v62_hunt_engine",
  "version": "62.1.0",
  "generated_at": "2026-06-28T02:41:04.575439+00:00",
  "total_hunts": 5,
  "active_campaigns": 16,
  "attack_paths": 5,
  "hunt_hypotheses": [
    {
      "hunt_id": "HUNT-405AAF62",
      "technique": "T1566",
      "hypothesis": "MFA bypass via adversary-in-the-middle phishing kit deployment",
      "priority": "HIGH",
      "confidence": 85,
      "status": "ACTIVE",
      "created_at": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "hunt_id": "HUNT-2CC5C805",
      "technique": "T1071.001",
      "hypothesis": "Botnet C2 infrastructure using domain generation algorithms to evade blocklists",
      "priority": "HIGH",
      "confidence": 82,
      "status": "ACTIVE",
      "created_at": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "hunt_id": "HUNT-93FF50B3",
      "technique": "T1547.001",
      "hypothesis": "Remote access trojan persistence via HKCU Run key and DLL side-loading",
      "priority": "HIGH",
      "confidence": 87,
      "status": "ACTIVE",
      "created_at": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "hunt_id": "HUNT-1DA2DC7C",
      "technique": "T1053.005",
      "hypothesis": "Nation-state long-term persistence via scheduled task and living-off-the-land techniques",
      "priority": "HIGH",
      "confidence": 78,
      "status": "ACTIVE",
      "created_at": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "hunt_id": "HUNT-1D7EF706",
      "technique": "T1195.002",
      "hypothesis": "Malicious package injection via compromised upstream supplier",
      "priority": "CRITICAL",
      "confidence": 89,
      "status": "ACTIVE",
      "created_at": "2026-06-28T02:41:04.575439+00:00"
    }
  ],
  "campaign_intelligence": [
    {
      "campaign_name": "vulnerability_exploitation_cluster_ops",
      "campaign_id": "CAMP-BC5AF389",
      "actors_involved": [
        "Vulnerability Exploitation Cluster"
      ],
      "incident_count": 64,
      "avg_risk": 7.8,
      "techniques_observed": [
        "T1190",
        "T1195",
        "T1547",
        "T1071",
        "T1041"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "cdb-apt-22_ops",
      "campaign_id": "CAMP-54A64124",
      "actors_involved": [
        "CDB-APT-22"
      ],
      "incident_count": 7,
      "avg_risk": 4.3,
      "techniques_observed": [
        "T1059",
        "T1566",
        "T1548",
        "T1190",
        "T1203",
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "cdb-fin-09_ops",
      "campaign_id": "CAMP-F92FD3FF",
      "actors_involved": [
        "CDB-FIN-09"
      ],
      "incident_count": 5,
      "avg_risk": 4.5,
      "techniques_observed": [
        "T1059",
        "T1566",
        "T1210",
        "T1190",
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "turla_(snake_/_venomous_bear)_ops",
      "campaign_id": "CAMP-FB58CABB",
      "actors_involved": [
        "Turla (Snake / Venomous Bear)"
      ],
      "incident_count": 4,
      "avg_risk": 6.9,
      "techniques_observed": [
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "cdb-fin-12_ops",
      "campaign_id": "CAMP-454F1DBC",
      "actors_involved": [
        "CDB-FIN-12"
      ],
      "incident_count": 4,
      "avg_risk": 6.7,
      "techniques_observed": [
        "T1059",
        "T1566",
        "T1548",
        "T1190",
        "T1195",
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "unattributed_ops",
      "campaign_id": "CAMP-E396CB85",
      "actors_involved": [
        "Unattributed"
      ],
      "incident_count": 4,
      "avg_risk": 5.5,
      "techniques_observed": [
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "phishing_threat_cluster_ops",
      "campaign_id": "CAMP-F0BA9B87",
      "actors_involved": [
        "Phishing Threat Cluster"
      ],
      "incident_count": 2,
      "avg_risk": 5.0,
      "techniques_observed": [
        "T1547",
        "T1566",
        "T1219"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "apt41_(double_dragon)_ops",
      "campaign_id": "CAMP-C4139B17",
      "actors_involved": [
        "APT41 (Double Dragon)"
      ],
      "incident_count": 1,
      "avg_risk": 7.2,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "muddywater_ops",
      "campaign_id": "CAMP-AF28EB8F",
      "actors_involved": [
        "MuddyWater"
      ],
      "incident_count": 1,
      "avg_risk": 7.5,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "lockbit_ops",
      "campaign_id": "CAMP-7278D717",
      "actors_involved": [
        "LockBit"
      ],
      "incident_count": 1,
      "avg_risk": 8.5,
      "techniques_observed": [
        "T1486"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "lazarus_group_ops",
      "campaign_id": "CAMP-7152C326",
      "actors_involved": [
        "Lazarus Group"
      ],
      "incident_count": 1,
      "avg_risk": 1.6,
      "techniques_observed": [
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "cdb-fin-11_ops",
      "campaign_id": "CAMP-1BBCE60E",
      "actors_involved": [
        "CDB-FIN-11"
      ],
      "incident_count": 1,
      "avg_risk": 1.0,
      "techniques_observed": [
        "T1078",
        "T1566",
        "T1203",
        "T1071"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "cdb-apt-28_ops",
      "campaign_id": "CAMP-F51A01C2",
      "actors_involved": [
        "CDB-APT-28"
      ],
      "incident_count": 1,
      "avg_risk": 5.5,
      "techniques_observed": [
        "T1059",
        "T1203",
        "T1190",
        "T1078",
        "T1055"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "criminal_ransomware_ecosystem_ops",
      "campaign_id": "CAMP-74383E18",
      "actors_involved": [
        "Criminal Ransomware Ecosystem"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "mustang_panda_ops",
      "campaign_id": "CAMP-259B9908",
      "actors_involved": [
        "Mustang Panda"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    },
    {
      "campaign_name": "lazarus_ops",
      "campaign_id": "CAMP-B1800A23",
      "actors_involved": [
        "Lazarus"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-28T02:41:04.575439+00:00"
    }
  ],
  "metrics": {
    "hypotheses_active": 5,
    "campaigns_tracked": 16,
    "avg_confidence": 84.2
  }
}