# data/threathunts/hunts.json
{
  "engine": "v62_hunt_engine",
  "version": "62.0.0",
  "generated_at": "2026-04-13T17:21:40.466543+00:00",
  "total_hunts": 5,
  "active_campaigns": 5,
  "attack_paths": 5,
  "hunt_hypotheses": [
    {
      "hunt_id": "HUNT-93FF50B3",
      "technique": "T1547.001",
      "hypothesis": "Remote access trojan persistence via HKCU Run key and DLL side-loading",
      "priority": "HIGH",
      "confidence": 87,
      "status": "ACTIVE",
      "created_at": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "hunt_id": "HUNT-1D7EF706",
      "technique": "T1195.002",
      "hypothesis": "Malicious package injection via compromised upstream supplier",
      "priority": "CRITICAL",
      "confidence": 89,
      "status": "ACTIVE",
      "created_at": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "hunt_id": "HUNT-1DA2DC7C",
      "technique": "T1053.005",
      "hypothesis": "Nation-state long-term persistence via scheduled task and living-off-the-land techniques",
      "priority": "HIGH",
      "confidence": 78,
      "status": "ACTIVE",
      "created_at": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "hunt_id": "HUNT-235AD79C",
      "technique": "T1078.004",
      "hypothesis": "Cloud environment takeover via stolen API keys with persistence via IAM role abuse",
      "priority": "HIGH",
      "confidence": 83,
      "status": "ACTIVE",
      "created_at": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "hunt_id": "HUNT-482A78C2",
      "technique": "T1041",
      "hypothesis": "Slow-and-low data exfiltration using encrypted C2 beaconing to cloud storage",
      "priority": "HIGH",
      "confidence": 80,
      "status": "ACTIVE",
      "created_at": "2026-04-13T17:21:40.466543+00:00"
    }
  ],
  "campaign_intelligence": [
    {
      "campaign_name": "lazarus_ops",
      "campaign_id": "CAMP-B1800A23",
      "actors_involved": [
        "Lazarus"
      ],
      "incident_count": 14,
      "avg_risk": 4.2,
      "techniques_observed": [
        "T1566",
        "T1041",
        "T1486",
        "T1567",
        "T1595",
        "T1547"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "campaign_name": "revil_ops",
      "campaign_id": "CAMP-BDAB2EA5",
      "actors_involved": [
        "REvil"
      ],
      "incident_count": 11,
      "avg_risk": 3.9,
      "techniques_observed": [
        "T1486",
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "campaign_name": "apt28_ops",
      "campaign_id": "CAMP-466B1115",
      "actors_involved": [
        "APT28"
      ],
      "incident_count": 9,
      "avg_risk": 4.8,
      "techniques_observed": [
        "T1566",
        "T1567",
        "T1595",
        "T1547",
        "T1036",
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "campaign_name": "apt41_ops",
      "campaign_id": "CAMP-4EFD834E",
      "actors_involved": [
        "APT41"
      ],
      "incident_count": 2,
      "avg_risk": 4.0,
      "techniques_observed": [
        "T1566",
        "T1071",
        "T1547",
        "T1078",
        "T1203",
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-04-13T17:21:40.466543+00:00"
    },
    {
      "campaign_name": "salt_typhoon_ops",
      "campaign_id": "CAMP-CFD1FF52",
      "actors_involved": [
        "Salt Typhoon"
      ],
      "incident_count": 1,
      "avg_risk": 2.9,
      "techniques_observed": [
        "T1190",
        "T1078"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-04-13T17:21:40.466543+00:00"
    }
  ],
  "metrics": {
    "hypotheses_active": 5,
    "campaigns_tracked": 5,
    "avg_confidence": 83.4
  }
}