# data/threathunts/hunts.json
{
  "engine": "v62_hunt_engine",
  "version": "62.1.0",
  "generated_at": "2026-06-26T20:48:47.758532+00:00",
  "total_hunts": 4,
  "active_campaigns": 11,
  "attack_paths": 4,
  "hunt_hypotheses": [
    {
      "hunt_id": "HUNT-93FF50B3",
      "technique": "T1547.001",
      "hypothesis": "Remote access trojan persistence via HKCU Run key and DLL side-loading",
      "priority": "HIGH",
      "confidence": 87,
      "status": "ACTIVE",
      "created_at": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "hunt_id": "HUNT-1DA2DC7C",
      "technique": "T1053.005",
      "hypothesis": "Nation-state long-term persistence via scheduled task and living-off-the-land techniques",
      "priority": "HIGH",
      "confidence": 78,
      "status": "ACTIVE",
      "created_at": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "hunt_id": "HUNT-405AAF62",
      "technique": "T1566",
      "hypothesis": "MFA bypass via adversary-in-the-middle phishing kit deployment",
      "priority": "HIGH",
      "confidence": 85,
      "status": "ACTIVE",
      "created_at": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "hunt_id": "HUNT-235AD79C",
      "technique": "T1078.004",
      "hypothesis": "Cloud environment takeover via stolen API keys with persistence via IAM role abuse",
      "priority": "HIGH",
      "confidence": 83,
      "status": "ACTIVE",
      "created_at": "2026-06-26T20:48:47.758532+00:00"
    }
  ],
  "campaign_intelligence": [
    {
      "campaign_name": "turla_(snake_/_venomous_bear)_ops",
      "campaign_id": "CAMP-FB58CABB",
      "actors_involved": [
        "Turla (Snake / Venomous Bear)"
      ],
      "incident_count": 6,
      "avg_risk": 6.2,
      "techniques_observed": [
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "lazarus_group_ops",
      "campaign_id": "CAMP-7152C326",
      "actors_involved": [
        "Lazarus Group"
      ],
      "incident_count": 5,
      "avg_risk": 4.0,
      "techniques_observed": [
        "T1190",
        "T1078",
        "T1059"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "cdb-fin-12_ops",
      "campaign_id": "CAMP-454F1DBC",
      "actors_involved": [
        "CDB-FIN-12"
      ],
      "incident_count": 3,
      "avg_risk": 7.1,
      "techniques_observed": [
        "T1548",
        "T1059",
        "T1041",
        "T1190",
        "T1078",
        "T1566"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "cdb-fin-09_ops",
      "campaign_id": "CAMP-F92FD3FF",
      "actors_involved": [
        "CDB-FIN-09"
      ],
      "incident_count": 3,
      "avg_risk": 3.9,
      "techniques_observed": [
        "T1190",
        "T1210",
        "T1078",
        "T1059"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "cdb-apt-22_ops",
      "campaign_id": "CAMP-54A64124",
      "actors_involved": [
        "CDB-APT-22"
      ],
      "incident_count": 3,
      "avg_risk": 1.6,
      "techniques_observed": [
        "T1190",
        "T1078",
        "T1203",
        "T1059"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "cdb-ran-03_ops",
      "campaign_id": "CAMP-78F2FC35",
      "actors_involved": [
        "CDB-RAN-03"
      ],
      "incident_count": 2,
      "avg_risk": 9.2,
      "techniques_observed": [
        "T1190",
        "T1059",
        "T1486"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "apt41_(double_dragon)_ops",
      "campaign_id": "CAMP-C4139B17",
      "actors_involved": [
        "APT41 (Double Dragon)"
      ],
      "incident_count": 1,
      "avg_risk": 7.2,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "muddywater_ops",
      "campaign_id": "CAMP-AF28EB8F",
      "actors_involved": [
        "MuddyWater"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "cdb-fin-11_ops",
      "campaign_id": "CAMP-1BBCE60E",
      "actors_involved": [
        "CDB-FIN-11"
      ],
      "incident_count": 1,
      "avg_risk": 1.0,
      "techniques_observed": [
        "T1071",
        "T1078",
        "T1203",
        "T1566"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "mustang_panda_ops",
      "campaign_id": "CAMP-259B9908",
      "actors_involved": [
        "Mustang Panda"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    },
    {
      "campaign_name": "lazarus_ops",
      "campaign_id": "CAMP-B1800A23",
      "actors_involved": [
        "Lazarus"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-26T20:48:47.758532+00:00"
    }
  ],
  "metrics": {
    "hypotheses_active": 4,
    "campaigns_tracked": 11,
    "avg_confidence": 83.2
  }
}