# data/threathunts/hunts.json
{
  "engine": "v62_hunt_engine",
  "version": "62.1.0",
  "generated_at": "2026-06-23T09:50:22.993677+00:00",
  "total_hunts": 5,
  "active_campaigns": 8,
  "attack_paths": 5,
  "hunt_hypotheses": [
    {
      "hunt_id": "HUNT-93FF50B3",
      "technique": "T1547.001",
      "hypothesis": "Remote access trojan persistence via HKCU Run key and DLL side-loading",
      "priority": "HIGH",
      "confidence": 87,
      "status": "ACTIVE",
      "created_at": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "hunt_id": "HUNT-235AD79C",
      "technique": "T1078.004",
      "hypothesis": "Cloud environment takeover via stolen API keys with persistence via IAM role abuse",
      "priority": "HIGH",
      "confidence": 83,
      "status": "ACTIVE",
      "created_at": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "hunt_id": "HUNT-1DA2DC7C",
      "technique": "T1053.005",
      "hypothesis": "Nation-state long-term persistence via scheduled task and living-off-the-land techniques",
      "priority": "HIGH",
      "confidence": 78,
      "status": "ACTIVE",
      "created_at": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "hunt_id": "HUNT-1D7EF706",
      "technique": "T1195.002",
      "hypothesis": "Malicious package injection via compromised upstream supplier",
      "priority": "CRITICAL",
      "confidence": 89,
      "status": "ACTIVE",
      "created_at": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "hunt_id": "HUNT-2CC5C805",
      "technique": "T1071.001",
      "hypothesis": "Botnet C2 infrastructure using domain generation algorithms to evade blocklists",
      "priority": "HIGH",
      "confidence": 82,
      "status": "ACTIVE",
      "created_at": "2026-06-23T09:50:22.993677+00:00"
    }
  ],
  "campaign_intelligence": [
    {
      "campaign_name": "turla_(snake_/_venomous_bear)_ops",
      "campaign_id": "CAMP-FB58CABB",
      "actors_involved": [
        "Turla (Snake / Venomous Bear)"
      ],
      "incident_count": 4,
      "avg_risk": 5.6,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "cdb-ran-03_ops",
      "campaign_id": "CAMP-78F2FC35",
      "actors_involved": [
        "CDB-RAN-03"
      ],
      "incident_count": 2,
      "avg_risk": 9.2,
      "techniques_observed": [
        "T1059",
        "T1486",
        "T1190"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "cdb-cyb-01_ops",
      "campaign_id": "CAMP-B46DB640",
      "actors_involved": [
        "CDB-CYB-01"
      ],
      "incident_count": 2,
      "avg_risk": 9.0,
      "techniques_observed": [
        "T1078",
        "T1041",
        "T1190",
        "T1059",
        "T1210",
        "T1566"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "lazarus_ops",
      "campaign_id": "CAMP-B1800A23",
      "actors_involved": [
        "Lazarus"
      ],
      "incident_count": 2,
      "avg_risk": 6.1,
      "techniques_observed": [
        "T1195"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "lockbit_ransomware_group_ops",
      "campaign_id": "CAMP-D144DBFF",
      "actors_involved": [
        "LockBit Ransomware Group"
      ],
      "incident_count": 1,
      "avg_risk": 8.6,
      "techniques_observed": [
        "T1486"
      ],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "muddywater_ops",
      "campaign_id": "CAMP-AF28EB8F",
      "actors_involved": [
        "MuddyWater"
      ],
      "incident_count": 1,
      "avg_risk": 7.6,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "lazarus_group_ops",
      "campaign_id": "CAMP-7152C326",
      "actors_involved": [
        "Lazarus Group"
      ],
      "incident_count": 1,
      "avg_risk": 7.1,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    },
    {
      "campaign_name": "mustang_panda_ops",
      "campaign_id": "CAMP-259B9908",
      "actors_involved": [
        "Mustang Panda"
      ],
      "incident_count": 1,
      "avg_risk": 5.0,
      "techniques_observed": [],
      "status": "ACTIVE",
      "last_activity": "2026-06-23T09:50:22.993677+00:00"
    }
  ],
  "metrics": {
    "hypotheses_active": 5,
    "campaigns_tracked": 8,
    "avg_confidence": 83.8
  }
}