# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--92768ee062e8",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580829+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--64707a756cb0",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in technology critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581080+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--3c59fd0458fb",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581113+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--c2e551178fde",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581196+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4869196137f",
    "hypothesis": "Ransomware deployment by UNC-CDB-INGEST using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580496+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7f5a0a4fdfc6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580565+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--402465847aaf",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in education critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580614+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--9c8fb3a6b930",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580668+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4d50dd80373",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580719+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--9de61e25f92a",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580762+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--102d10232dd4",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580795+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--0553b0cd369e",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580862+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--e0f633120622",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580894+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--8d03c7d43123",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580928+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--478a0f87b935",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580961+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--71e63495ff79",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in technology critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.580990+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--af870e23a24e",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581021+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--21da28385514",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581052+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--8cd7b8794e57",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581145+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--8e23d0487dc5",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-27T19:52:36.581228+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  }
]