# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--3e4698890842",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1595",
      "T1195",
      "T1190",
      "T1078",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.019891+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-07"
    ]
  },
  {
    "hunt_id": "hunt--2d9459f7dd9f",
    "hypothesis": "Threat actor CDB-FIN-09 conducting credential harvesting campaign using T1566 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1566",
      "T1190",
      "T1078",
      "T1059.001",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020048+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--374a5b7fc4f1",
    "hypothesis": "Nation-state actor CDB-FIN-09 establishing persistence via T1195 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1195",
      "T1199",
      "T1078",
      "T1203",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020106+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--f28f542fb087",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020178+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--fb60b5d56cfe",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1059",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020210+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--af2d9e2134de",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1059",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020241+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--c17ef7db7e8f",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1059",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020271+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--139b169597ab",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1595",
      "T1195",
      "T1199",
      "T1190",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020341+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--6744bb4f6c63",
    "hypothesis": "Nation-state actor CDB-RAN-01 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1566.001",
      "T1199",
      "T1190",
      "T1476"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020378+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-RAN-01"
    ]
  },
  {
    "hunt_id": "hunt--2f751ac917b5",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059",
      "T1176",
      "T1068",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020412+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--854778b5d925",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1595 in retail critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1203",
      "T1555",
      "T1567"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020447+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--8935c1f26140",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1566.002",
      "T1059",
      "T1547",
      "T1027",
      "T1620"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020481+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--d828e8030b09",
    "hypothesis": "Nation-state actor CDB-FIN-12 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1566",
      "T1195",
      "T1078",
      "T1190"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020512+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-12"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1203 in energy critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1567"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020632+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--0538f6232469",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1195",
      "T1190",
      "T1078",
      "T1547.001",
      "T1412"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.019962+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--1c1008a1b1d0",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020006+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--8b42a38a5370",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020142+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--6c349338cf2a",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1566",
      "T1195",
      "T1199",
      "T1190",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020310+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-RAN-03"
    ]
  },
  {
    "hunt_id": "hunt--57a5b9adf90c",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020549+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--b8cb62a20327",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-03T07:44:39.020603+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  }
]