# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--2f874ebcc8f0",
    "hypothesis": "Nation-state actor CDB-FIN-07 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1566",
      "T1190",
      "T1199",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378684+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-07"
    ]
  },
  {
    "hunt_id": "hunt--59887fdf7cf4",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1566",
      "T1199",
      "T1190",
      "T1203",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378749+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-29"
    ]
  },
  {
    "hunt_id": "hunt--f550c00d6df0",
    "hypothesis": "Nation-state actor CDB-APT-29 establishing persistence via T1566 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1199",
      "T1078",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378799+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-29"
    ]
  },
  {
    "hunt_id": "hunt--5add7e1a6a34",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1583.001 in financial critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1583.001",
      "T1566",
      "T1566.002",
      "T1199",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378849+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--87e82317c1ca",
    "hypothesis": "Nation-state actor CDB-APT-28 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1199",
      "T1078",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378934+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-28"
    ]
  },
  {
    "hunt_id": "hunt--c341cba925e5",
    "hypothesis": "Nation-state actor CDB-MOB-01 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1566.002",
      "T1199",
      "T1078",
      "T1475"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378998+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-MOB-01"
    ]
  },
  {
    "hunt_id": "hunt--08ea248ac856",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-99 targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1203",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.379078+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1203",
      "T1542",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.379190+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--106a1c629e0a",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.378898+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--ee07f881045a",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1595",
      "T1542",
      "T1556",
      "T1046",
      "T1071"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.379036+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--da6cbf36dbbf",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1203 against energy infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1555",
      "T1111"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.379115+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--844b0eb4343a",
    "hypothesis": "Threat actor UNC-CDB-SWARM-AI conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-11T13:40:37.379249+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-SWARM-AI"
    ]
  }
]