# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--5add7e1a6a34",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1595 in government critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1078.001",
      "T1078",
      "T1203",
      "T1204.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770152+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--7ec473ab5faa",
    "hypothesis": "Threat actor CDB-FIN-12 conducting credential harvesting campaign using T1566 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1566",
      "T1566.002",
      "T1204.001",
      "T1036",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770205+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-12"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1195 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1195",
      "T1203",
      "T1204.001",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770250+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--87e82317c1ca",
    "hypothesis": "Nation-state actor CDB-APT-28 establishing persistence via T1593 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1593",
      "T1566",
      "T1566.002",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770425+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-28"
    ]
  },
  {
    "hunt_id": "hunt--106a1c629e0a",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1059.001",
      "T1059",
      "T1547",
      "T1547.001",
      "T1053.005"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770460+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--08ea248ac856",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-99 targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1078",
      "T1203",
      "T1059",
      "T1547",
      "T1505.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770510+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--844b0eb4343a",
    "hypothesis": "Threat actor UNC-CDB-SWARM-AI conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770090+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-SWARM-AI"
    ]
  },
  {
    "hunt_id": "hunt--8b8419dad637",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770299+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--da6cbf36dbbf",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1190 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1190",
      "T1078",
      "T1203",
      "T1204.001",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770339+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--c0f441fc333e",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-08T18:52:30.770384+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  }
]