# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--5add7e1a6a34",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1190",
      "T1203",
      "T1547",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178545+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--d72a5ae8d57b",
    "hypothesis": "Ransomware deployment by CDB-FIN-11 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1595",
      "T1566",
      "T1078",
      "T1199",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178666+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-11"
    ]
  },
  {
    "hunt_id": "hunt--374a5b7fc4f1",
    "hypothesis": "Nation-state actor CDB-FIN-09 establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1204.002",
      "T1078",
      "T1453",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178753+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--106a1c629e0a",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1566",
      "T1078",
      "T1203",
      "T1059",
      "T1547"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178945+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--82131097823c",
    "hypothesis": "Nation-state actor CDB-FIN-11 establishing persistence via T1199 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1199",
      "T1203",
      "T1059",
      "T1068",
      "T1021"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178978+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-11"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1203",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179019+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--ee07f881045a",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1595",
      "T1078",
      "T1059",
      "T1555",
      "T1213"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179054+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--3473e0ddff54",
    "hypothesis": "Ransomware deployment by CDB-FIN-09 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1566",
      "T1195",
      "T1078",
      "T1555",
      "T1539"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179151+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--08ea248ac856",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-99 targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1078",
      "T1203",
      "T1059",
      "T1547",
      "T1505.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179184+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--286b2a0f8294",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by CDB-FIN-11 targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1595",
      "T1566",
      "T1078",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179213+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-11"
    ]
  },
  {
    "hunt_id": "hunt--c341cba925e5",
    "hypothesis": "Nation-state actor CDB-MOB-01 establishing persistence via T1195 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1195",
      "T1199",
      "T1190",
      "T1078",
      "T1476"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179245+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-MOB-01"
    ]
  },
  {
    "hunt_id": "hunt--507f1c742aa7",
    "hypothesis": "Ransomware deployment by CDB-RAN-02 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1190",
      "T1530",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179299+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-RAN-02"
    ]
  },
  {
    "hunt_id": "hunt--f1a41e1deb57",
    "hypothesis": "Nation-state actor CDB-APT-41 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1078",
      "T1476",
      "T1203",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179357+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-41"
    ]
  },
  {
    "hunt_id": "hunt--d828e8030b09",
    "hypothesis": "Nation-state actor CDB-FIN-12 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1566",
      "T1195",
      "T1078",
      "T1199"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.179420+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-12"
    ]
  },
  {
    "hunt_id": "hunt--9d7d529c0b94",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178624+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--1c6850772dce",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178706+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--844b0eb4343a",
    "hypothesis": "Threat actor UNC-CDB-SWARM-AI conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178791+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-SWARM-AI"
    ]
  },
  {
    "hunt_id": "hunt--da6cbf36dbbf",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059",
      "T1555",
      "T1213"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178832+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--71e3681bbdd2",
    "hypothesis": "Ransomware deployment by CDB-FIN-07 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1078",
      "T1190",
      "T1203",
      "T1555",
      "T1021.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-18T19:20:07.178883+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-FIN-07"
    ]
  }
]