# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--5add7e1a6a34",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1036",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710595+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1195",
      "T1203",
      "T1059",
      "T1027"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710656+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--374a5b7fc4f1",
    "hypothesis": "Nation-state actor CDB-FIN-09 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1199",
      "T1190",
      "T1078",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710709+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--9f468172a1bc",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1566",
      "T1195",
      "T1078",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710755+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--f550c00d6df0",
    "hypothesis": "Nation-state actor CDB-APT-29 establishing persistence via T1583.001 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1583.001",
      "T1566",
      "T1566.002",
      "T1195",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710793+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-29"
    ]
  },
  {
    "hunt_id": "hunt--5ba59a0fc635",
    "hypothesis": "Ransomware deployment by CDB-APT-22 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1190",
      "T1203",
      "T1555",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710869+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--f1a41e1deb57",
    "hypothesis": "Nation-state actor CDB-APT-41 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1078",
      "T1476",
      "T1203",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710985+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-41"
    ]
  },
  {
    "hunt_id": "hunt--da6cbf36dbbf",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1195 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1195"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710832+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--f7c9fdefa824",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.710918+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--844b0eb4343a",
    "hypothesis": "Threat actor UNC-CDB-SWARM-AI conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-10T03:12:04.711024+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-SWARM-AI"
    ]
  }
]