# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--82131097823c",
    "hypothesis": "Nation-state actor CDB-FIN-11 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1195",
      "T1190",
      "T1078",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189385+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-11"
    ]
  },
  {
    "hunt_id": "hunt--f969ac314fcb",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by CDB-IR-02 targeting internet-facing government systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1595",
      "T1566",
      "T1199",
      "T1078",
      "T1476"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189447+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-IR-02"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1203 in government critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189618+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--2f751ac917b5",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1566 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1059.001",
      "T1059",
      "T1547",
      "T1547.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189654+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--dd2981598fb7",
    "hypothesis": "Nation-state actor CDB-CYB-01 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1195",
      "T1078",
      "T1059",
      "T1547"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189721+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--4837b83e194f",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1203",
      "T1059",
      "T1071.004",
      "T1041",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189766+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--7ec473ab5faa",
    "hypothesis": "Threat actor CDB-FIN-12 conducting credential harvesting campaign using T1595 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1595",
      "T1566",
      "T1566.002",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189916+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-12"
    ]
  },
  {
    "hunt_id": "hunt--acfa97872878",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1059",
      "T1176"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.190095+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--9f11fc8037b1",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189500+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--70fbd33233da",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189544+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--153e6db87edc",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189584+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--d8bdf6b8252b",
    "hypothesis": "Ransomware deployment by CDB-RAN-01 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1203",
      "T1059",
      "T1547",
      "T1547.001",
      "T1053.005"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189687+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-RAN-01"
    ]
  },
  {
    "hunt_id": "hunt--b48684c8ea3b",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189804+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--31d2bf741f72",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189840+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--e3b8c5c736be",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528",
      "T1078",
      "T1190",
      "T1203",
      "T1204.001",
      "T1530"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189874+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--293ce9314c54",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.189965+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--db37decd8b8a",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.190025+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--a93d9751162a",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.190061+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--fdcd620b773b",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1190 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1190",
      "T1059",
      "T1213"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.190126+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--91d4a7c62663",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-31T14:04:32.190157+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  }
]