# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1595 in energy critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1190",
      "T1078",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808427+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--08ea248ac856",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-99 targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068",
      "T1078",
      "T1203",
      "T1059",
      "T1547",
      "T1505.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808489+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--5add7e1a6a34",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1195 in retail critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1195",
      "T1078",
      "T1203",
      "T1036",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808611+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--3473e0ddff54",
    "hypothesis": "Ransomware deployment by CDB-FIN-09 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1566",
      "T1195",
      "T1078",
      "T1555",
      "T1539"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808771+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--da6cbf36dbbf",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1566 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1566",
      "T1190",
      "T1078",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808810+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--374a5b7fc4f1",
    "hypothesis": "Nation-state actor CDB-FIN-09 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1199",
      "T1190",
      "T1078",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808860+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-09"
    ]
  },
  {
    "hunt_id": "hunt--9f468172a1bc",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1566",
      "T1195",
      "T1078",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808892+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--f550c00d6df0",
    "hypothesis": "Nation-state actor CDB-APT-29 establishing persistence via T1583.001 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1583.001",
      "T1566",
      "T1566.002",
      "T1195",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808922+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-29"
    ]
  },
  {
    "hunt_id": "hunt--5ba59a0fc635",
    "hypothesis": "Ransomware deployment by CDB-APT-22 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1190",
      "T1203",
      "T1555",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808958+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--ed9f82856249",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1203 in retail critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1036",
      "T1555",
      "T1539",
      "T1111"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808535+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--844b0eb4343a",
    "hypothesis": "Threat actor UNC-CDB-SWARM-AI conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808575+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-SWARM-AI"
    ]
  },
  {
    "hunt_id": "hunt--26ba228d1130",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808671+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--ca71db3b6afe",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808707+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--8599bc12efc6",
    "hypothesis": "Threat actor UNC-CDB-99 conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.808739+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  },
  {
    "hunt_id": "hunt--f7c9fdefa824",
    "hypothesis": "Nation-state actor UNC-CDB-99 establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-09T19:12:29.809004+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-99"
    ]
  }
]