# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--dd2981598fb7",
    "hypothesis": "Nation-state actor CDB-CYB-01 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1195",
      "T1078",
      "T1190",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784828+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--2f751ac917b5",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1547",
      "T1021",
      "T1041",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784978+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--2821c9f40884",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1190 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1190",
      "T1078",
      "T1203",
      "T1547",
      "T1542"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785115+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--0fbb9c1bde2b",
    "hypothesis": "Nation-state actor CDB-IR-02 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785181+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-IR-02"
    ]
  },
  {
    "hunt_id": "hunt--d72a5ae8d57b",
    "hypothesis": "Ransomware deployment by CDB-FIN-11 using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001",
      "T1595",
      "T1566",
      "T1078",
      "T1199",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785279+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-FIN-11"
    ]
  },
  {
    "hunt_id": "hunt--4eed2ec2e349",
    "hypothesis": "Threat actor CDB-CYB-01 conducting credential harvesting campaign using T1566 against financial infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1566",
      "T1195",
      "T1199",
      "T1078",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785484+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--dc1ba747944b",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784668+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--4f11e9ee9201",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784730+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--53a7204efb03",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against retail infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784782+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--03eae28b1b81",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784872+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--664a70f9ae65",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059",
      "T1071"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784909+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--17a133a35431",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1190 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1190",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.784944+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--1a47abea2850",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785011+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--84796d675db3",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785043+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--c19a755ceaf4",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785080+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--327b6d49fb07",
    "hypothesis": "Threat actor CDB-IR-02 conducting credential harvesting campaign using T1583.001 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1583.001",
      "T1566",
      "T1566.002",
      "T1190",
      "T1078"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785229+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-IR-02"
    ]
  },
  {
    "hunt_id": "hunt--d4b383c63275",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785323+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--702c4ad39704",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against financial infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785352+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--88c6bcbfe030",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785385+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--c8ce599df5cc",
    "hypothesis": "Threat actor CDB-RU-01 conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1547",
      "T1542",
      "T1021",
      "T1412"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-29T19:04:23.785428+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "CDB-RU-01"
    ]
  }
]