# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--8dba6a054cdb",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386892+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--dd848fbb8432",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386925+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--d35ca90fc0c2",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386952+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4869196137f",
    "hypothesis": "Ransomware deployment by UNC-CDB-INGEST using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386266+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--402465847aaf",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386332+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--5ae5ef8d431a",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against financial infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386386+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--a7f72e395cd8",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386443+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--256217d34c87",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386475+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--c735b1284968",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386504+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--a6db0523080e",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386532+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--64f6f5c3c593",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386556+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--16d25c7f9372",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in retail critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386592+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--d302c12dbecc",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386621+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--74503090a8fa",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in retail critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386649+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7f5a0a4fdfc6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against retail infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386690+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--80647df21ae5",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386752+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--bd3f12733945",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386798+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--add84a987fee",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386829+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--cfb6274a2905",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386861+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--55ec4af97007",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-06T08:58:51.386981+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  }
]