# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--ebeac9e0d35c",
    "hypothesis": "Threat actor UNK conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.674507+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--fca96f028cae",
    "hypothesis": "Nation-state actor UNK establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.674604+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--a6f0ed695574",
    "hypothesis": "Threat actor UNK conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.674698+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--9ca535f1d6a2",
    "hypothesis": "Threat actor UNK conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.674882+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--eb65cf52c9f1",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.674946+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--c220cbb6ad0b",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNK targeting internet-facing government systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.675000+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--5a6e31f0c457",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.675039+00:00",
    "priority": "HIGH",
    "actor_tags": []
  },
  {
    "hunt_id": "hunt--ba70704b2747",
    "hypothesis": "Ransomware deployment by UNK using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-14T08:06:02.675076+00:00",
    "priority": "HIGH",
    "actor_tags": []
  }
]