# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--402465847aaf",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141509+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7f5a0a4fdfc6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141580+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--4602fef0745b",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-INGEST targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141650+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--71e76692cab4",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141689+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--d9ab1a2360c3",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141730+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--07a114a4b213",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141760+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--d2000c1b75ad",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141788+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--6ef7c406bb43",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141817+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--527e9abbe447",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against retail infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141849+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--83c0b56dc162",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141903+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--bd4c0d07291c",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141935+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4ddc93c1821",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting enterprise SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141964+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4869196137f",
    "hypothesis": "Ransomware deployment by UNC-CDB-INGEST using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.141990+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--14d032488883",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.142062+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--461b9f8a98bf",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.142175+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--c497a19f0007",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-04-29T14:33:22.142206+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  }
]