# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--2f751ac917b5",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203",
      "T1547",
      "T1542",
      "T1071",
      "T1486"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.053927+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--433a2c3683ab",
    "hypothesis": "Nation-state actor CDB-RU-02 establishing persistence via T1566 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1566",
      "T1059.001",
      "T1059",
      "T1547",
      "T1053.005"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.053990+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-RU-02"
    ]
  },
  {
    "hunt_id": "hunt--2600bba1230c",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1595",
      "T1078",
      "T1203",
      "T1059",
      "T1547"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054043+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CN-01"
    ]
  },
  {
    "hunt_id": "hunt--39d7ede3015c",
    "hypothesis": "Nation-state actor CDB-APT-22 establishing persistence via T1595 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1595",
      "T1078",
      "T1190",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054093+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--8935c1f26140",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1195",
      "T1078",
      "T1059",
      "T1547",
      "T1547.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054200+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-APT-22"
    ]
  },
  {
    "hunt_id": "hunt--dd2981598fb7",
    "hypothesis": "Nation-state actor CDB-CYB-01 establishing persistence via T1195 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1195",
      "T1190",
      "T1078",
      "T1547",
      "T1555"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054731+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "CDB-CYB-01"
    ]
  },
  {
    "hunt_id": "hunt--46988a832197",
    "hypothesis": "Nation-state actor UNC-UNKNOWN establishing persistence via T1203 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003",
      "T1203"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054783+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--e735f662bb95",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059",
      "T1213"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054132+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--150c66289266",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054238+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--2dc10014eb32",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1190 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1190",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054272+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--491ce3bae9fa",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054307+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--4a8662ccd7fd",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054340+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--57a5b9adf90c",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1595 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1595",
      "T1195",
      "T1190",
      "T1203",
      "T1547"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054373+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--fabe26244c49",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054408+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--0b4294a715af",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054440+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--9ed667a74c4b",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1203 against technology infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1203",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054474+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--4837b83e194f",
    "hypothesis": "Covert data exfiltration campaign using encrypted channels to adversary-controlled infrastructure",
    "data_sources": [
      "Network Connection",
      "DNS Query",
      "Proxy Logs",
      "Cloud API Logs"
    ],
    "techniques": [
      "T1041",
      "T1048",
      "T1567",
      "T1071.001",
      "T1078",
      "T1059",
      "T1204.001",
      "T1071"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054516+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--75f343606dd6",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054570+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--8621c504b416",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059",
      "T1499"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054601+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  },
  {
    "hunt_id": "hunt--148d3f6ed9fa",
    "hypothesis": "Threat actor UNC-UNKNOWN conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110",
      "T1059"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-03-26T14:01:02.054658+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-UNKNOWN"
    ]
  }
]