# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--0cba16d6938e",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310161+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--61097b674e46",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310186+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--23a5ac5954cb",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309715+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--5b6b703d6113",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309763+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--33ed4b88c9ce",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309798+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--acef1c502b3d",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309825+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--bc63d3e9c5c1",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309852+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--4602fef0745b",
    "hypothesis": "Active exploitation of unpatched vulnerability N/A by UNC-CDB-INGEST targeting internet-facing enterprise systems",
    "data_sources": [
      "Web Server Logs",
      "IDS/IPS Alerts",
      "Process Creation",
      "File Creation"
    ],
    "techniques": [
      "T1190",
      "T1203",
      "T1068"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309874+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7f5a0a4fdfc6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against retail infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309895+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4ddc93c1821",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting technology SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309928+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--2dd949b48c0b",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309967+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--22df6038f5b2",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.309989+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--c6003f1c23d6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310020+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7a4f4c90f8c7",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting retail SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310043+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--79f5c8c268db",
    "hypothesis": "Cloud infrastructure compromise via stolen API keys or OAuth token abuse targeting retail SaaS environments",
    "data_sources": [
      "Cloud Audit Logs",
      "API Call Logs",
      "Identity Provider Logs",
      "OAuth Token Events"
    ],
    "techniques": [
      "T1078.004",
      "T1550.001",
      "T1528"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310066+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--225dbc2deaa1",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310091+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--566311affe81",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310111+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--28d3829af30c",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310136+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b372fe469bf1",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310229+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b7c5e525fbad",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-06-07T14:23:33.310256+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  }
]