# data/nexus/threat_hunts.json
[
  {
    "hunt_id": "hunt--ed9a3eefbac1",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.458183+00:00",
    "priority": "CRITICAL",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b4869196137f",
    "hypothesis": "Ransomware deployment by UNC-CDB-INGEST using living-off-the-land binaries for lateral movement before encryption",
    "data_sources": [
      "Process Creation",
      "File Modification",
      "Service Creation",
      "Network Share Access"
    ],
    "techniques": [
      "T1486",
      "T1490",
      "T1021.002",
      "T1059.001"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457512+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--402465847aaf",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in enterprise critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457589+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7f5a0a4fdfc6",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457635+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--c5a9807df6b7",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457684+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--f8607d2d587a",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457724+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--e03d9c3206e2",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457803+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--8497aabc9d6d",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457839+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--7dcbc7033479",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457871+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--9e75556f540a",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in critical infrastructure critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457904+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--b9675774c37c",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against critical infrastructure infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457938+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--e33e71a51156",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.457970+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--4493b454083f",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.458079+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--201a9fe2ee49",
    "hypothesis": "Nation-state actor UNC-CDB-INGEST establishing persistence via T1059 in technology critical infrastructure",
    "data_sources": [
      "Registry Modification",
      "Scheduled Task",
      "Service Creation",
      "WMI Activity"
    ],
    "techniques": [
      "T1547.001",
      "T1053.005",
      "T1543.003",
      "T1546.003"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.458114+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--3b8a137cd0b5",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.458150+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  },
  {
    "hunt_id": "hunt--a73f19ac8249",
    "hypothesis": "Threat actor UNC-CDB-INGEST conducting credential harvesting campaign using T1059 against enterprise infrastructure",
    "data_sources": [
      "Authentication Logs",
      "Process Access",
      "Windows Event Log",
      "Network Connection"
    ],
    "techniques": [
      "T1003",
      "T1003.001",
      "T1003.006",
      "T1110"
    ],
    "status": "PLANNED",
    "findings": [],
    "detection_rules": [],
    "created_at": "2026-05-24T19:43:26.458218+00:00",
    "priority": "HIGH",
    "actor_tags": [
      "UNC-CDB-INGEST"
    ]
  }
]